# meta-mct/meta-s5549/recipes-phosphor/gpio/power-state/SetPowerGoodPropertyOff.service
[Unit]
Description=Set Power Good Property Off

[Service]
Type=oneshot
ExecStart=/usr/sbin/setPowerProperty.sh off
