# meta-mct/meta-s5549/recipes-phosphor/mct-dcmi/mct-dcmi-power/mct-dcmi-power.service
[Unit]
Description=MiTAC DCMI Power
Requires=xyz.openbmc_project.psusensor.service
After=xyz.openbmc_project.psusensor.service

[Service]
Restart=always
ExecStart=/usr/bin/env mct-dcmi-power
SyslogIdentifier=mct-dcmi-power
TimeoutSec=infinity

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
