# meta-mct/meta-s7106/recipes-phosphor/gpio/caterr-update/xyz.openbmc_project.caterr.service
[Unit]
Description=assert CATERR 

[Service]
Type=oneshot
ExecStart=/usr/bin/env caterr_update.sh assert
SyslogIdentifier=caterr_assert

