# meta-mct/meta-s8030/recipes-phosphor/mct-early-post/mct-early-post-bin/set-vga-source-poweroff.service
[Unit]
Description=Set VGA source when DC power off

[Service]
Restart=always
ExecStart=/usr/bin/env set-vga-source-poweroff.sh
SyslogIdentifier=set-vga-source-poweroff
TimeoutSec=infinity

[Install]
WantedBy=sysinit.target
