# meta-mct/meta-f5i/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.prochot1.service
[Unit]
Description=assert Prochot1 

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochot_update.sh assert CPU1_State 
SyslogIdentifier=prochot1_assert

