# meta-mct/meta-f5i/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.prochot0.service
[Unit]
Description=assert Prochot0 

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochot_update.sh assert CPU0_State 
SyslogIdentifier=prochot0_assert

