# meta-mct/meta-f5i/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.caterr_deassert.service
[Unit]
Description=deassert CATERR 

[Service]
Type=oneshot
ExecStart=/usr/bin/env caterr_update.sh deassert
SyslogIdentifier=caterr_deassert

