# meta-mct/meta-s5549/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.prochot0_deassert.service
[Unit]
Description=deassert Prochot0 
StartLimitIntervalSec=6
StartLimitBurst=1

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochotThermtripHandler deassert CPU0_State prochot
SyslogIdentifier=prochot0_deassert
RemainAfterExit=true