# meta-mct/meta-s5549/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.thermtrip.service
[Unit]
Description=assert Thermtrip 

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochotThermtripHandler assert CPU_State thermtrip
SyslogIdentifier=thermtrip_assert

