# meta-mct/meta-s5549/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.caterr.service
[Unit]
Description=assert CATERR 

[Service]
Type=oneshot
ExecStart=/usr/bin/caterrHandler
SyslogIdentifier=caterr_assert

