# meta-mct/meta-common/recipes-phosphor/ipmbbridged/files/xyz.openbmc_project.Ipmi.Channel.Ipmb@.service
[Unit]
Description=IPMB bridged daemon

After=phosphor-ipmi-host.service
[Service]
Restart=always
ExecStart=/usr/bin/env ipmbbridged --p %I
SyslogIdentifier=ipmbbridged
TimeoutSec=infinity
