# meta-mct/meta-s8030/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.thermtrip.service
[Unit]
Description=assert Thermtrip 

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochot_thermtrip_update.sh assert CPU_State thermtrip 
SyslogIdentifier=thermtrip_assert

