# meta-mct/meta-s8036/recipes-phosphor/s8036-default-action/files/sync-time-initial.service
[Unit]
Description=Timesync Initialize

[Service]
Restart=no
ExecStart=/usr/sbin/sync-time-initial.py
SyslogIdentifier=sync-time-initial
TimeoutSec=infinity

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
