# meta-mct/meta-s8036/recipes-phosphor/gpio/phosphor-gpio-monitor/SetPowerGoodPropertyOn.service
[Unit]
Description=Set Power Good Property On

[Service]
Type=oneshot
ExecStart=/usr/sbin/setPowerProperty.sh on

