# meta-mct/meta-s5549/recipes-phosphor/mct-default-action/files/disable-jtag.service
[Unit]
Description=Set the jtag driver default disabled

[Service]
Type=simple
ExecStart=/usr/bin/env disable-jtag.sh
SyslogIdentifier=disable-jtag

[Install]
WantedBy=sysinit.target

