# meta-mct/meta-f5i/recipes-phosphor/gpio/pwr-button/xyz.openbmc_project.pwrbutton.service
[Unit]
Description=Log Power buttong press 

[Service]
Type=oneshot
ExecStart=/usr/bin/env pwrbutton_update.sh
SyslogIdentifier=pwrbutton_press

